# launchd/com.brain.execution-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brain.execution-api</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bard/Code/brain/api/execution_log_api.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/bard/Code/brain</string>
    
    <key>StandardOutPath</key>
    <string>/Users/bard/Code/brain/logs/api.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bard/Code/brain/logs/api_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
